# files/supervisor/supervisord.conf
[supervisord]
nodaemon=true
autostart=true
autorestart=true
user=root

[program:dockergen]
command=docker-gen -watch -only-exposed -notify "pkill -sighup tor" /app/torrc.tmpl /etc/torrc
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tor]
command=/usr/bin/tor -f /etc/torrc
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
